# misc/redis-server.conf
[program:redis-server]
command=/path/to/redis-server <set your option>
autostart=true
autorestart=true
user=<your username>
redirect_stderr=true
stdout_logfile=NONE
priority=1
